# estimate/backend-cpp/packaging/aorms-estimate-backend.service
[Unit]
Description=AORMS Estimate — C++ + SQLite estimating backend
Documentation=https://aorms.in
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
# Overridable settings live in /etc/default/aorms-estimate-backend.
EnvironmentFile=-/etc/default/aorms-estimate-backend
# Default DB under the service state dir (created by StateDirectory below).
Environment=ESTIMATE_DB=/var/lib/aorms-estimate/aorms-estimate.db
Environment=ESTIMATE_HOST=127.0.0.1
Environment=ESTIMATE_PORT=8787
ExecStart=/usr/bin/aorms-estimate-backend
Restart=on-failure
RestartSec=2

# Run under an isolated, auto-provisioned system user with a private state dir.
DynamicUser=yes
StateDirectory=aorms-estimate
ProtectSystem=strict
ProtectHome=yes
PrivateTmp=yes
NoNewPrivileges=yes

[Install]
WantedBy=multi-user.target
